# common/tmp/etc/systemd/system/getty@tty1.service.d/skip-username.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --login-options '-p -- usrj' --noclear --nonewline --skip-login - $TERM
Environment=XDG_SESSION_TYPE=x11
