# common/Documents/tmp_for_install/etc/NetworkManager/dispatcher.d/10-timezone
#!/bin/sh

# [ "$2" = up ] && timedatectl set-timezone "$(curl --fail https://ipapi.co/timezone)"
