# common/Documents/tmp_for_install/etc/NetworkManager/dispatcher.d/11-ntpdate
#!/bin/sh

[ "$2" = up ] && systemctl start ntpdate.service
